# service/sidekiq/run
#!/usr/bin/env bash

set -e

sdir="$(cd -P . && pwd -P)"

cd ../..
exec 2>&1
exec chpst -e "$sdir/env" bash -c 'exec bundle exec sidekiq -r ./lib/worker.rb'
